# sms-gateway.service
[Unit]
Description=SMS Gateway Cloudflare Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/run_gateway_services.py
WorkingDirectory=/root
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
